# packaging/systemd/immich-frame-kiosk.service
[Unit]
Description=Immich Frame Chromium kiosk
Wants=immich-frame.service
After=immich-frame.service systemd-user-sessions.service

[Service]
Type=simple
User=immich-frame-kiosk
Group=immich-frame-kiosk
EnvironmentFile=/etc/immich-frame/kiosk.env
WorkingDirectory=/var/lib/immich-frame/kiosk
ExecStart=/usr/bin/startx /usr/local/lib/immich-frame/start-kiosk.sh -- -nocursor
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
